# omri-dashboard.service
[Unit]
Description=Omri Association Dashboard
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/Omri-Association
Environment=PATH=/path/to/Omri-Association/venv/bin
ExecStart=/path/to/Omri-Association/venv/bin/streamlit run dashboard.py --server.port 8501 --server.address 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
